# templates/searchindex/root/etc/init/create-overview-user.conf.erb
description "Create the 'overview' user"
author "Overview Project"

start on starting overview-search-index
start on starting overview-message-broker

pre-start script
  grep overview /etc/group || groupadd overview -g 998
  grep overview /etc/passwd || useradd overview -c "Overview Servers" -d /opt/overview -g 998 --system --shell /bin/false -u 998
  chown overview:overview /opt/overview
end script
