# templates/searchindex/root/etc/init/overview-search-index.conf.erb
description "Overview search index"
author "Overview Project"


start on runlevel [2345]
stop on runlevel [016]

setuid overview
setgid overview
chdir /opt/overview
console output

respawn


exec /usr/bin/java \
  -cp "<%= managed_code_path %>/search-index/*" \
  -Des.path.home=/var/lib/elasticsearch \
  -Xms2g \
  -Xmx2g \
  -Xss256k \
  -XX:+UseParNewGC \
  -XX:+UseConcMarkSweepGC \
  -XX:CMSInitiatingOccupancyFraction=75 \
  -XX:+UseCMSInitiatingOccupancyOnly \
  -Djava.awt.headless=true \
  -Delasticsearch \
  -Des.foreground=yes \
  org.elasticsearch.bootstrap.Bootstrap \
  2>&1 | logger -p local1.notice
