# templates/production/worker/root/etc/init/overview-worker.conf.erb
description "Overview worker"
author "Overview Project"

start on runlevel [2345]
stop on runlevel [016]

setuid overview
setgid overview
chdir /opt/overview
console output

respawn

env SMTP_USERNAME="<%= secrets['smtp-username'] %>"
env SMTP_PASSWORD="<%= secrets['smtp-password'] %>"
env SMTP_HOST="<%= smtp_host %>"
env SMTP_PORT="<%= smtp_port %>"
env SMTP_SSL="<%= smtp_ssl && 'true' || 'false' %>"
env LOGGER_SMTP_FROM="<%= logger_smtp_from %>"
env LOGGER_SMTP_TO="<%= logger_smtp_to %>"
env LOGGER_SMTP_SUBJECT_PREFIX="<%= logger_worker_smtp_subject_prefix %>"

<%# TODO: remove database_url from the commandline %>
exec /usr/bin/java \
       -cp "<%= managed_code_path %>/main/current/*" \
       -Ddatasource.default.driver=org.postgresql.Driver \
       -Ddatasource.default.url="<%= secrets['database-url'] %>" \
       -Dlogback.configurationFile=workerprodlog.xml \
       JobHandler \
       2>&1 | logger -p local1.notice
