# templates/web/root/etc/init/haproxy.conf.erb
description "HAProxy SSL endpoint for Overview Project"
author "Overview Project"

pre-start script
  grep haproxy /etc/group || groupadd haproxy -g 997
  grep haproxy /etc/passwd || useradd haproxy -c "HAProxy server" -d /tmp -g 997 --system --shell /bin/false -u 997
end script

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /opt/haproxy/sbin/haproxy -f <%= config_path %>/haproxy.cfg
